# resources/syncthing_install.list
sudo apt install curl apt-transport-https
sudo echo "deb https://apt.syncthing.net/ syncthing stable" | sudo tee /etc/apt/sources.list.d/syncthing.list
curl -s https://syncthing.net/release-key.txt | sudo apt-key add -
sudo apt-get update
sudo apt-get install syncthing
sudo cp resource/syncthing@.service /etc/systemd/system/.
sudo systemctl daemon-reload
sudo systemctl start syncthing@$USER
sudo systemctl enable syncthing@$USER
sudo systemctl status syncthing@$USER
